# main/airootfs/etc/systemd/system/power-profiles-daemon.service.d/profile.conf
[Service]
ExecStart=
ExecStart=/usr/bin/powerprofilesctl set performance

[Install]
WantedBy=
WantedBy=multi-user.target
